# systemd/rclone-azureblob@.service
# /etc/systemd/system/rclone-azureblob@.service
#
# Azure Blob Storage 전용 rclone mount 인스턴스 서비스
#
# 사용법:
#   # 1. 설정 파일 생성
#   sudo install -m 640 /dev/null /etc/rclone/mounts/myblob.conf
#   sudo tee /etc/rclone/mounts/myblob.conf << 'EOF'
#   REMOTE=myremote:mycontainer
#   MOUNTPOINT=/mnt/azureblob
#   EXTRA_ARGS=--allow-other --vfs-cache-mode writes --vfs-cache-max-size 2G
#   EOF
#
#   # 2. (선택) Azure 인증 환경 변수 별도 파일로 관리
#   sudo install -m 600 /dev/null /etc/rclone/azure-env.conf
#   sudo tee -a /etc/rclone/azure-env.conf << 'EOF'
#   AZURE_TENANT_ID=...
#   AZURE_CLIENT_ID=...
#   AZURE_CLIENT_SECRET=...
#   EOF
#
#   # 3. 서비스 시작
#   sudo mkdir -p /mnt/azureblob
#   sudo systemctl start  rclone-azureblob@myblob.service
#   sudo systemctl enable rclone-azureblob@myblob.service

[Unit]
Description=rclone Azure Blob mount: %i
Documentation=https://rclone.org/azureblob/
After=network-online.target
Wants=network-online.target
# 주의: ConditionPathIsDirectory 는 RuntimeDirectory 생성 전에 평가되므로
#       /run/rclone-azureblob/%i 를 조건으로 쓰면 항상 실패함 → 제거

[Service]
Type=notify

# ── 마운트 인스턴스 설정 (필수) ──────────────────────────────────────────────
# /etc/rclone/mounts/<INSTANCE>.conf 에서 REMOTE, MOUNTPOINT, EXTRA_ARGS 읽음
EnvironmentFile=/etc/rclone/mounts/%i.conf

# ── Azure 인증 환경 변수 (선택: env_auth 방식 사용 시) ───────────────────────
# 파일이 없으면 무시됨 (-를 붙이면 파일 없어도 오류 없음)
EnvironmentFile=-/etc/rclone/azure-env.conf

# ── 런타임 준비 ───────────────────────────────────────────────────────────────
# RuntimeDirectory가 먼저 생성된 후 ExecStartPre가 실행됨
ExecStartPre=/bin/mkdir -p ${MOUNTPOINT}

# ── rclone mount 실행 ─────────────────────────────────────────────────────────
# /usr/bin/rclone: deb 설치 시 실제 경로
# /usr/local/bin/rclone: tarball 설치 또는 심볼릭 링크 경로 (하위 호환)
ExecStart=/usr/bin/rclone mount \
    --config /etc/rclone/rclone.conf \
    --log-level INFO \
    --log-file /var/log/rclone/azureblob-%i.log \
    \
    --vfs-cache-mode writes \
    --vfs-cache-max-size 2G \
    --vfs-cache-max-age 1h \
    --vfs-write-back 5s \
    --buffer-size 256M \
    \
    --dir-cache-time 30s \
    --poll-interval 15s \
    \
    --timeout 30s \
    --contimeout 15s \
    --low-level-retries 10 \
    --retries 3 \
    --retries-sleep 1s \
    \
    ${EXTRA_ARGS:-} \
    ${REMOTE} ${MOUNTPOINT}

# ── 정지 처리 ─────────────────────────────────────────────────────────────────
ExecStop=/bin/fusermount3 -u ${MOUNTPOINT}

# ── 재시작 정책 ───────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=10s
StartLimitInterval=120s
StartLimitBurst=3

# ── 리소스 ────────────────────────────────────────────────────────────────────
LimitNOFILE=1048576
TimeoutStartSec=60
TimeoutStopSec=30

RuntimeDirectory=rclone-azureblob/%i
RuntimeDirectoryMode=0755
LogsDirectory=rclone
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target
