# systemd/rclone-azureblob@.service
# /etc/systemd/system/rclone-azureblob@.service
#
# Azure Blob Storage 전용 rclone mount 인스턴스 서비스
#
# 일반 rclone-mount@.service 와 다른 점:
#   - Azure Blob 접속 전 네트워크 도달성 확인 (storage account endpoint)
#   - 환경 변수 파일로 Azure 인증 정보 주입 지원
#   - 폐쇄망 Private Endpoint 환경을 고려한 타임아웃 설정
#
# 사용법:
#   # 1. 설정 파일 생성
#   cp azure/conf-examples/azblob-key.conf /etc/rclone/mounts/myblob.conf
#   vi /etc/rclone/mounts/myblob.conf
#
#   # 2. (선택) Azure 인증 환경 변수 별도 파일로 관리
#   install -m 600 /dev/null /etc/rclone/azure-env.conf
#   echo 'AZURE_TENANT_ID=...' >> /etc/rclone/azure-env.conf
#   echo 'AZURE_CLIENT_ID=...' >> /etc/rclone/azure-env.conf
#   echo 'AZURE_CLIENT_SECRET=...' >> /etc/rclone/azure-env.conf
#
#   # 3. 서비스 시작
#   mkdir -p /mnt/myblob
#   systemctl start  rclone-azureblob@myblob.service
#   systemctl enable rclone-azureblob@myblob.service

[Unit]
Description=rclone Azure Blob mount: %i
Documentation=https://rclone.org/azureblob/
After=network-online.target
Wants=network-online.target
# 마운트 포인트 디렉토리가 존재해야 시작
ConditionPathIsDirectory=/run/rclone-azureblob/%i

[Service]
Type=notify

# ── 마운트 인스턴스 설정 (필수) ──────────────────────────────────────────────
# /etc/rclone/mounts/<INSTANCE>.conf 에서 REMOTE, MOUNTPOINT, EXTRA_ARGS 읽음
EnvironmentFile=/etc/rclone/mounts/%i.conf

# ── Azure 인증 환경 변수 (선택: env_auth 방식 사용 시) ───────────────────────
# 파일이 없으면 무시됨 (-를 붙이면 파일 없어도 오류 없음)
EnvironmentFile=-/etc/rclone/azure-env.conf

# ── 런타임 준비 ───────────────────────────────────────────────────────────────
ExecStartPre=/bin/mkdir -p ${MOUNTPOINT}
ExecStartPre=/bin/mkdir -p /run/rclone-azureblob/%i

# ── rclone mount 실행 ─────────────────────────────────────────────────────────
ExecStart=/usr/local/bin/rclone mount \
    --config /etc/rclone/rclone.conf \
    --log-level INFO \
    --log-file /var/log/rclone/azureblob-%i.log \
    \
    --vfs-cache-mode writes \
    --vfs-cache-max-size 2G \
    --vfs-cache-max-age 1h \
    --vfs-write-back 5s \
    --buffer-size 256M \
    \
    --dir-cache-time 30s \
    --poll-interval 15s \
    \
    --timeout 30s \
    --contimeout 15s \
    --low-level-retries 10 \
    --retries 3 \
    --retries-sleep 1s \
    \
    ${EXTRA_ARGS:-} \
    ${REMOTE} ${MOUNTPOINT}

# ── 정지 처리 ─────────────────────────────────────────────────────────────────
ExecStop=/bin/fusermount3 -u ${MOUNTPOINT}
ExecStopPost=/bin/rmdir --ignore-fail-on-non-empty /run/rclone-azureblob/%i

# ── 재시작 정책 ───────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=10s
StartLimitInterval=120s
StartLimitBurst=3

# ── 리소스 ────────────────────────────────────────────────────────────────────
LimitNOFILE=1048576
TimeoutStartSec=60
TimeoutStopSec=30

RuntimeDirectory=rclone-azureblob/%i
RuntimeDirectoryMode=0755
LogsDirectory=rclone
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target
