# systemd/rclone-mount@.service
# /etc/systemd/system/rclone-mount@.service
#
# 범용 rclone mount 인스턴스 서비스 (template unit)
# Azure Blob 전용: rclone-azureblob@.service 사용 권장
#
# 사용법:
#   systemctl start  rclone-mount@<INSTANCE>.service
#   systemctl enable rclone-mount@<INSTANCE>.service
#
# 설정 파일 형식 (/etc/rclone/mounts/<INSTANCE>.conf):
#   REMOTE=myremote:bucket/path
#   MOUNTPOINT=/mnt/s3data
#   EXTRA_ARGS=--vfs-cache-mode writes --allow-other

[Unit]
Description=rclone mount: %i
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
EnvironmentFile=/etc/rclone/mounts/%i.conf

# RuntimeDirectory가 먼저 생성된 후 ExecStartPre가 실행됨
ExecStartPre=/bin/mkdir -p ${MOUNTPOINT}

# /usr/bin/rclone: deb 설치 경로 (tarball은 /usr/local/bin/rclone 심볼릭 링크)
# 주의: systemd ExecStart 는 shell 기본값 확장(${VAR:-})을 지원하지 않으므로
#       다중 인자 옵션은 $EXTRA_ARGS 로 whitespace splitting 되게 전달한다.
ExecStart=/usr/bin/rclone mount \
    --config /etc/rclone/rclone.conf \
    --log-level INFO \
    --log-file /var/log/rclone/%i.log \
    --vfs-cache-mode writes \
    --vfs-cache-max-size 1G \
    --vfs-cache-max-age 1h \
    --buffer-size 256M \
    --dir-cache-time 30s \
    --poll-interval 15s \
    $EXTRA_ARGS \
    ${REMOTE} ${MOUNTPOINT}

ExecStop=/bin/fusermount3 -u ${MOUNTPOINT}

Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

LimitNOFILE=1048576
TimeoutStopSec=30

RuntimeDirectory=rclone-mount/%i
RuntimeDirectoryMode=0755
LogsDirectory=rclone
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target
